# gcore-on-fault.service
[Service]
Type=exec
ExecStart=/usr/local/libexec/gcore-on-fault.sh

[Install]
WantedBy=multi-user.target
